# roles/jmeter-slave/files/jmeter.conf
[program:jmeter]
command=/opt/run-server.sh
user=root
autostart=true
autorestart=true
stderr_logfile=/var/log/jmeter.log
stdout_logfile=/var/log/jmeter.log
